# dokuwiki/monitor.cron
* * * * * /usr/bin/php /var/www/html/doku.php
